# rootfs/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:ase-dataserver]
command=/ase-start-dataserver.sh

[program:ase-backupserver]
command=/ase-start-backupserver.sh
